# deploy/systemd/linaholod.service
# systemd-юнит для Next.js (next start).
# Кладётся в /etc/systemd/system/linaholod.service.
# Требует Node 20+ системно (через NodeSource → /usr/bin/npm).

[Unit]
Description=linaholod.ru (Next.js)
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/var/www/linaholod
Environment=NODE_ENV=production
Environment=PORT=3000
# .env.local в WorkingDirectory Next подхватывает сам (в т.ч. серверные TELEGRAM_*).
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
# Мягкая остановка: дать завершить запросы и after()-колбэки.
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
